# launchd/com.chenie.iwatch-voice.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chenie.iwatch-voice</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/chenie/Documents/project/iwatch-voice/dist/src/cli.js</string>
    <string>watch</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/chenie/Documents/project/iwatch-voice</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/chenie/Documents/project/iwatch-voice/runtime/launchd.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/chenie/Documents/project/iwatch-voice/runtime/launchd.stderr.log</string>
</dict>
</plist>
